# cookbooks/hello_world/templates/onetime.service.erb
[Unit]
Description=My One-Time Service

[Service]
ExecStart=/usr/bin/chef-client

[Install]
WantedBy=multi-user.target



